# archive/nodecontroller/nc-baseos/rootfs/usr/share/upstart/sessions/upstart-file-bridge.conf
# upstart-file-bridge - Bridge file events into upstart
#
# This helper daemon receives inotify(7) events and
# emits equivalent Upstart events.

description	"Bridge file events into upstart"

emits file

start on started dbus
stop on session-end

expect daemon
respawn

exec upstart-file-bridge --daemon --user
