# archive/nodecontroller/nc-baseos/rootfs/etc/init/monit.conf
 # This is an upstart script to keep monit running.
 # To install disable the old way of doing things:
 #
 #   /etc/init.d/monit stop && update-rc.d -f monit remove
 #
 # then put this script here:
 #
 #   /etc/init/monit.conf
 #
 # and reload upstart configuration:
 #
 #   initctl reload-configuration
 #
 # You can manually start and stop monit like this:
 # 
 # start monit
 # stop monit
 #

description "Monit service manager"

limit core unlimited unlimited

start on launch-monit
stop on runlevel [06]

expect daemon
respawn

exec /usr/bin/monit -c /etc/monit/monitrc
#exec /bin/bash /usr/bin/launch_monit.sh

pre-stop exec /usr/bin/monit -c /etc/monit/monitrc quit

