# archive/nodecontroller/nc-baseos/rootfs/etc/profile.d/custom.source

LOCKSTEP_FILE=/var/log/self-test/lockstep.log
export LOCKSTEP_FILE

SELF_TEST_APP_DIR=/root/self-test/apps/
export SELF_TEST_APP_DIR

WAGGLE_PID=/var/run/launch_waggle.pid
export WAGGLE_PID

SHUTDOWN_FILE=/var/proper_shutdown.waggle
export SHUTDOWN_FILE

SELF_TEST_LOG_DIR=/var/log/stress-test/
export SELF_TEST_LOG_DIR

