# nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/apt.conffiles
/etc/logrotate.d/apt
/etc/kernel/postinst.d/apt-auto-removal
/etc/cron.daily/apt
/etc/apt/apt.conf.d/01autoremove
/etc/apt/apt.conf.d/20changelog
