# nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/monit.list
/.
/usr
/usr/share
/usr/share/bug
/usr/share/bug/monit
/usr/share/doc
/usr/share/doc/monit
/usr/share/doc/monit/NEWS.Debian.gz
/usr/share/doc/monit/examples
/usr/share/doc/monit/examples/monit.upstart
/usr/share/doc/monit/examples/monit.service
/usr/share/doc/monit/README
/usr/share/doc/monit/copyright
/usr/share/doc/monit/README.Debian
/usr/share/doc/monit/changelog.Debian.gz
/usr/share/doc/monit/PLATFORMS.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/monit.1.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/monit
/usr/bin
/usr/bin/monit
/var
/var/lib
/var/lib/monit
/var/lib/monit/events
/etc
/etc/pam.d
/etc/pam.d/monit
/etc/default
/etc/default/monit
/etc/init.d
/etc/init.d/monit
/etc/logrotate.d
/etc/logrotate.d/monit
/etc/monit
/etc/monit/conf.d
/etc/monit/templates
/etc/monit/templates/rootstrict
/etc/monit/templates/rootrc
/etc/monit/templates/rootbin
/etc/monit/monitrc.d
/etc/monit/monitrc.d/postfix
/etc/monit/monitrc.d/smartmontools
/etc/monit/monitrc.d/apache2
/etc/monit/monitrc.d/memcached
/etc/monit/monitrc.d/openntpd
/etc/monit/monitrc.d/pdns-recursor
/etc/monit/monitrc.d/mdadm
/etc/monit/monitrc.d/mysql
/etc/monit/monitrc.d/nginx
/etc/monit/monitrc.d/rsyslog
/etc/monit/monitrc.d/cron
/etc/monit/monitrc.d/snmpd
/etc/monit/monitrc.d/acpid
/etc/monit/monitrc.d/at
/etc/monit/monitrc.d/openssh-server
/etc/monit/monitrc
