# nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/upstart.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/upstart
/usr/share/doc/upstart/NEWS.gz
/usr/share/doc/upstart/HACKING.gz
/usr/share/doc/upstart/TODO.gz
/usr/share/doc/upstart/copyright
/usr/share/doc/upstart/changelog.Debian.gz
/usr/share/doc/upstart/AUTHORS
/usr/share/doc/upstart/README.Debian.gz
/usr/share/upstart
/usr/share/upstart/sessions
/usr/share/upstart/sessions/xsession-init.conf
/usr/share/upstart/sessions/logrotate.conf
/usr/share/upstart/sessions/upstart-dbus-session-bridge.conf
/usr/share/upstart/sessions/upstart-file-bridge.conf
/usr/share/upstart/sessions/re-exec.conf
/usr/share/upstart/sessions/upstart-event-bridge.conf
/usr/share/upstart/sessions/upstart-dbus-system-bridge.conf
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/upstart.py
/usr/share/man
/usr/share/man/man5/inittab.5.gz
/usr/share/man/man5/init.5.gz
/usr/share/man/man7/starting.7.gz
/usr/share/man/man7/session-end.7.gz
/usr/share/man/man7/power-status-changed.7.gz
/usr/share/man/man7/upstart-events.7.gz
/usr/share/man/man7/stopping.7.gz
/usr/share/man/man7/runlevel.7.gz
/usr/share/man/man7/socket-event.7.gz
/usr/share/man/man7/stopped.7.gz
/usr/share/man/man7/keyboard-request.7.gz
/usr/share/man/man7/file-event.7.gz
/usr/share/man/man7/control-alt-delete.7.gz
/usr/share/man/man7/dbus-event.7.gz
/usr/share/man/man7/startup.7.gz
/usr/share/man/man7/started.7.gz
/usr/share/man/man8/initctl2dot.8.gz
/usr/share/man/man8/telinit.8.gz
/usr/share/man/man8/shutdown.8.gz
/usr/share/man/man8/upstart-dbus-bridge.8.gz
/usr/share/man/man8/upstart-local-bridge.8.gz
/usr/share/man/man8/runlevel.8.gz
/usr/share/man/man8/upstart-file-bridge.8.gz
/usr/share/man/man8/init-checkconf.8.gz
/usr/share/man/man8/initctl.8.gz
/usr/share/man/man8/reboot.8.gz
/usr/share/man/man8/init.8.gz
/usr/share/man/man8/upstart-event-bridge.8.gz
/usr/share/man/man8/upstart-socket-bridge.8.gz
/usr/share/man/man8/upstart-udev-bridge.8.gz
/usr/lib
/usr/lib/upstart
/usr/lib/upstart/migrate-inittab.pl
/usr/bin
/usr/bin/init-checkconf
/usr/bin/initctl2dot
/var
/var/log
/var/log/upstart
/sbin
/sbin/shutdown
/sbin/init
/sbin/initctl
/sbin/upstart-socket-bridge
/sbin/reboot
/sbin/upstart-event-bridge
/sbin/upstart-local-bridge
/sbin/upstart-file-bridge
/sbin/upstart-udev-bridge
/sbin/telinit
/sbin/upstart-dbus-bridge
/sbin/runlevel
/etc
/etc/init
/etc/init/upstart-socket-bridge.conf
/etc/init/upstart-udev-bridge.conf
/etc/init/failsafe.conf
/etc/init/container-detect.conf
/etc/init/tty1.conf
/etc/init/rc-sysinit.conf
/etc/init/shutdown.conf
/etc/init/tty6.conf
/etc/init/rcS.conf
/etc/init/tty5.conf
/etc/init/upstart-file-bridge.conf
/etc/init/tty2.conf
/etc/init/control-alt-delete.conf
/etc/init/wait-for-state.conf
/etc/init/flush-early-job-log.conf
/etc/init/tty3.conf
/etc/init/console.conf
/etc/init/rc.conf
/etc/init/tty4.conf
/etc/cron.daily
/etc/cron.daily/upstart
/etc/bash_completion.d
/etc/bash_completion.d/upstart
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/Upstart.conf
/etc/upstart-xsessions
/etc/X11
/etc/X11/Xsession.d
/etc/X11/Xsession.d/99upstart
/etc/X11/Xsession.d/00upstart
/etc/logrotate.d
/etc/logrotate.d/upstart
/lib
/lib/init
/lib/init/upstart-job
/lib/init/apparmor-profile-load
/bin
/bin/running-in-container
/usr/share/man/man7/net-device-added.7.gz
/usr/share/man/man7/drm-device-added.7.gz
/usr/share/man/man7/net-device-removed.7.gz
/usr/share/man/man7/upstart.7.gz
/usr/share/man/man7/graphics-device-added.7.gz
/usr/share/man/man8/poweroff.8.gz
/usr/share/man/man8/reload.8.gz
/usr/share/man/man8/halt.8.gz
/usr/share/man/man8/restart.8.gz
/usr/share/man/man8/status.8.gz
/usr/share/man/man8/stop.8.gz
/usr/share/man/man8/start.8.gz
/sbin/poweroff
/sbin/status
/sbin/stop
/sbin/start
/sbin/restart
/sbin/halt
/sbin/reload
