# core/fault_resilience/node/sw/Odroid_U3_HWRev0.5/docs/designConsiderations/data/improper_shutdown.conf
# /etc/init/improper_shutdown_detect.conf

description "Detect improper shutdown and auto-correct system"
author " "

env SHUTDOWN_FILE=/var/proper_shutdown.waggle

start on runlevel [2345]
stop on runlevel [06]


pre-start script
if [ -e "${SHUTDOWN_FILE}" ]
then
    logger -is -t "$UPSTART_JOB" "`date`: File exists at start. Normal start up"
    exec rm "${SHUTDOWN_FILE}"
else
    logger -is -t "$UPSTART_JOB" "`date`: File does not exist. Improper shutdown detected"
fi 
end script


# Things to run on (re)start
#post-start script 
#end script

# Things to run before shutdown
post-stop script

if [ -e "${SHUTDOWN_FILE}" ]
then
    logger -is -t "$UPSTART" "`date`: Lock file existed during post-stop"
else
    exec touch "${SHUTDOWN_FILE}"
fi
end script
