# archive/nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/keyboard-configuration.conffiles
/etc/init.d/console-setup
/etc/init/console-setup.conf
