# archive/nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/keyboard-configuration.list
/.
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/85-keyboard-configuration.rules
/lib/udev/console-setup-tty
/usr
/usr/share
/usr/share/console-setup
/usr/share/console-setup/keyboard
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/panic
/usr/share/initramfs-tools/scripts/panic/console_setup
/usr/share/initramfs-tools/scripts/init-top
/usr/share/initramfs-tools/scripts/init-top/console_setup
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/console_setup
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_console-setup.py
/usr/share/doc
/usr/share/doc/keyboard-configuration
/usr/share/doc/keyboard-configuration/README.Debian
/usr/share/doc/keyboard-configuration/copyright
/usr/share/doc/keyboard-configuration/examples
/usr/share/doc/keyboard-configuration/changelog.gz
/etc
/etc/init.d
/etc/init.d/console-setup
/etc/init
/etc/init/console-setup.conf
/usr/share/doc/keyboard-configuration/examples/keyboard
