# archive/nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/kmod.conffiles
/etc/init/kmod.conf
/etc/modprobe.d/blacklist-ath_pci.conf
/etc/modprobe.d/mlx4.conf
/etc/modprobe.d/blacklist-framebuffer.conf
/etc/modprobe.d/blacklist-firewire.conf
/etc/modprobe.d/iwlwifi.conf
/etc/modprobe.d/blacklist.conf
/etc/modprobe.d/blacklist-watchdog.conf
/etc/modprobe.d/blacklist-rare-network.conf
/etc/init.d/kmod
/etc/depmod.d/ubuntu.conf
