# archive/nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/monit.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/monit" ] || [ -e "/etc/init/monit.conf" ]; then
	invoke-rc.d monit stop || exit $?
fi
# End automatically added section
