# archive/nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/passwd.conffiles
/etc/cron.daily/passwd
/etc/default/useradd
/etc/init/passwd.conf
/etc/pam.d/passwd
/etc/pam.d/newusers
/etc/pam.d/chfn
/etc/pam.d/chpasswd
/etc/pam.d/chsh
