# archive/nodecontroller/nc-baseos/rootfs/var/lib/dpkg/info/sudo.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/sudo.service
/etc
/etc/init.d
/etc/init.d/sudo
/etc/pam.d
/etc/pam.d/sudo
/etc/sudoers
/etc/sudoers.d
/etc/sudoers.d/README
/usr
/usr/lib
/usr/lib/sudo
/usr/lib/sudo/sudo_noexec.so
/usr/lib/sudo/sudoers.so
/usr/lib/sudo/system_group.so
/usr/lib/sudo/group_file.so
/usr/lib/sudo/sesh
/usr/share
/usr/share/doc
/usr/share/doc/sudo
/usr/share/doc/sudo/README
/usr/share/doc/sudo/UPGRADE.gz
/usr/share/doc/sudo/NEWS.Debian.gz
/usr/share/doc/sudo/TROUBLESHOOTING.gz
/usr/share/doc/sudo/README.Debian
/usr/share/doc/sudo/copyright
/usr/share/doc/sudo/CONTRIBUTORS
/usr/share/doc/sudo/HISTORY
/usr/share/doc/sudo/OPTIONS
/usr/share/doc/sudo/changelog.Debian.gz
/usr/share/doc/sudo/NEWS.gz
/usr/share/doc/sudo/examples
/usr/share/doc/sudo/examples/sample.pam
/usr/share/doc/sudo/examples/sample.sudoers.gz
/usr/share/doc/sudo/examples/sample.syslog.conf
/usr/share/doc/sudo/examples/sample.sudo.conf
/usr/share/man
/usr/share/man/man5/sudo.conf.5.gz
/usr/share/man/man5/sudoers.5.gz
/usr/share/man/man8/visudo.8.gz
/usr/share/man/man8/sudoreplay.8.gz
/usr/share/man/man8/sudo.8.gz
/usr/share/man/man8/sudo_root.8.gz
/usr/share/man/man8/sudo_plugin.8.gz
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_sudo.py
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/sudo
/usr/sbin
/usr/sbin/visudo
/usr/bin
/usr/bin/sudo
/usr/bin/sudoreplay
/usr/include
/usr/include/sudo_plugin.h
/var
/var/lib
/var/lib/sudo
/usr/share/man/man8/sudoedit.8.gz
/usr/bin/sudoedit
